# .config/systemd/user/gammastep.service
[Unit]
After=graphical-session.target
[Service]
ExecStart=gammastep
[Install]
WantedBy=graphical-session.target
