# .config/systemd/user/gammastep.service
[Service]
ExecStart=gammastep -O 4000K
[Install]
WantedBy=default.target
