# .config/systemd/user/sway-wallpaper.timer
[Timer]
OnActiveSec=5s
OnUnitActiveSec=5m

[Install]
WantedBy=sway-session.target
